# .overlay/etc/systemd/system/dropcache.service
[Unit]
Description=Free up host memory

[Service]
Type=oneshot
ExecStart=sh -c 'echo 3 > /proc/sys/vm/drop_caches'
