# .overlay/etc/systemd/system/dropcache.timer
[Timer]
OnBootSec=0min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
